# lib/systemd/senselink-server.service
[Unit]
Description=SenseLink USB/IP server (forwards DualSense to remote host)
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
ExecStartPre=/sbin/modprobe usbip-core
ExecStartPre=/sbin/modprobe usbip-host
ExecStart=__BINDIR__/usbipd -D
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
